# deploy/tencent/shanyuan-backup.service
[Unit]
Description=Shanyuan Pavilion daily data backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/shanyuan/deploy/tencent/backup.sh