# deploy/tencent/shanyuan-backup.timer
[Unit]
Description=Run Shanyuan Pavilion backup every day

[Timer]
OnCalendar=*-*-* 03:25:00 UTC
Persistent=true

[Install]
WantedBy=timers.target